# deploy/hermes-line-gateway.service.example
[Unit]
Description=Hermes LINE WordPress Gateway
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/hermes-line-wordpress-agent
EnvironmentFile=-/etc/hermes-line-wordpress-agent/line-gateway.env
Environment=APP_DIR=/opt/hermes-line-wordpress-agent
Environment=ENV_FILE=/etc/hermes-line-wordpress-agent/line-gateway.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/hermes-line-wordpress-agent/deploy/run-line-gateway.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
